# spirrow-magickit-mcp.service
# Systemd service file for Spirrow Magickit MCP Server
# Copy to /etc/systemd/system/ and modify paths as needed
#
# Usage:
#   sudo cp spirrow-magickit-mcp.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable spirrow-magickit-mcp
#   sudo systemctl start spirrow-magickit-mcp

[Unit]
Description=Spirrow Magickit MCP Server
Documentation=https://github.com/SpirrowGames/spirrow-magickit
After=network.target
After=spirrow-magickit.service
Wants=spirrow-cognilens.service spirrow-prismind.service spirrow-lexora.service

[Service]
Type=simple
# TODO: Change to your user/group
User=YOUR_USER
Group=YOUR_GROUP
# TODO: Change to your installation path
WorkingDirectory=/path/to/spirrow-magickit
ExecStart=/path/to/spirrow-magickit/venv/bin/python -m magickit.mcp_server
Restart=always
RestartSec=5

# Environment
# TODO: Change to your installation path
Environment=PYTHONPATH=/path/to/spirrow-magickit/src
Environment=MAGICKIT_LOG_LEVEL=INFO

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
# TODO: Change to your installation path
ReadWritePaths=/path/to/spirrow-magickit/data
PrivateTmp=true

[Install]
WantedBy=multi-user.target
